# debian/example.cron
#!/bin/sh
# Cron.daily script for debmonitor, should be run with the debmonitor system user

/usr/bin/debmonitor
